# twtfwd.timer
[Unit]
Description=Twitter Bot

[Timer]
OnCalendar=*:0/5
Unit=twtfwd
 
[Install]
WantedBy=timers.target
