# twtfwd.service
[Unit]
Description=twtfwd

[Service]
Type=oneshot
ExecStart=/usr/bin/twtfwd
User=root
Group=root
