# build/twtfwd.service
[Unit]
Description=twtfwd

[Service]
Type=simple
ExecStart=/usr/bin/twtfwd
User=root
Group=root

[Install]
WantedBy = multi-user.target
